# claude-ops/scripts/systemd/pocket-email-bridge.service
[Unit]
Description=Pocket Email Bridge — outbound send + inbound reply parsing via Gmail
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
Group=__USER__
WorkingDirectory=__HOME__
# Doppler service token (same scoped token as pocket-watcher.service). Injects
# GOG_KEYRING_PASSWORD + GOG_ACCOUNT into the bridge's env. Without these, gog
# returns exit 2: "no TTY available for keyring file backend password prompt"
# or "missing --account (or set GOG_ACCOUNT)". See pocket-watcher.service
# operator setup notes.
LoadCredential=doppler-token:/etc/pocket/doppler-token
Environment=HOME=__HOME__
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=POCKET_STATE_DIR=__HOME__/.claude/state/pocket
Environment=GOG_BIN=/usr/local/bin/gog
ExecStart=/bin/bash -c 'DOPPLER_TOKEN=$(cat $CREDENTIALS_DIRECTORY/doppler-token) /usr/local/bin/doppler run --silent -- __HOME__/.venv-pocket/bin/python3 __PLUGIN_ROOT__/scripts/ops-pocket-email-bridge.py'
StandardOutput=append:__HOME__/.claude/state/pocket/email-bridge.log
StandardError=append:__HOME__/.claude/state/pocket/email-bridge-stderr.log
