# claude-ops/scripts/systemd/pocket-watcher.timer
[Unit]
Description=Run pocket-watcher every 5 minutes
Requires=pocket-watcher.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
