# claude-ops/scripts/systemd/pocket-watcher.service
[Unit]
Description=Pocket Watcher — polls Pocket AI MCP for new recordings
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__HOME__
# Secrets (POCKET_API_KEY, etc.) — see scripts/systemd/secrets.env.example.
# This file MUST be created by the operator at /etc/pocket/secrets.env
# with mode 600 owned by root. It is gitignored.
EnvironmentFile=/etc/pocket/secrets.env
Environment=HOME=__HOME__
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=POCKET_STATE_DIR=__HOME__/.claude/state/pocket
Environment=POCKET_MEMORY_DIR=__HOME__/.claude/projects/-__USER__/memory
Environment=POCKET_INDEX_FILE=__HOME__/.claude/projects/-__USER__/memory/MEMORY.md
Environment=POCKET_LOOKBACK_HOURS=24
Environment=GIGA_SYNC=1
Environment=GIGA_MIND=global
Environment=GIGA_CONSOLIDATE=1
ExecStart=__HOME__/.venv-pocket/bin/python3 __HOME__/claude-ops/claude-ops/scripts/ops-cron-pocket-watcher.py
StandardOutput=append:__HOME__/.claude/state/pocket/watcher-stdout.log
StandardError=append:__HOME__/.claude/state/pocket/watcher-stderr.log
