# claude-ops/scripts/systemd/pocket-watcher.service
[Unit]
Description=Pocket Watcher — polls Pocket AI MCP for new recordings
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
Group=__USER__
WorkingDirectory=__HOME__
# Doppler service token, scoped read-only to the operator's secrets config
# (e.g. claude-ops/dev). The token file is root-owned mode 600; the
# LoadCredential directive exposes it as $CREDENTIALS_DIRECTORY/doppler-token
# only to this service unit. doppler run then injects POCKET_API_KEY into
# the watcher's env at runtime — secrets never sit on disk in plaintext or
# in `systemctl cat` output.
#
# Operator setup:
#   sudo install -d -m 0750 -o root -g root /etc/pocket
#   echo "dp.st.dev.YOUR_TOKEN" | sudo tee /etc/pocket/doppler-token >/dev/null
#   sudo chmod 600 /etc/pocket/doppler-token && sudo chown root:root /etc/pocket/doppler-token
#   sudo systemctl daemon-reload && sudo systemctl restart pocket-watcher.timer
LoadCredential=doppler-token:/etc/pocket/doppler-token
Environment=HOME=__HOME__
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=POCKET_STATE_DIR=__HOME__/.claude/state/pocket
Environment=POCKET_MEMORY_DIR=__HOME__/.claude/projects/-__USER__/memory
Environment=POCKET_INDEX_FILE=__HOME__/.claude/projects/-__USER__/memory/MEMORY.md
Environment=POCKET_LOOKBACK_HOURS=24
Environment=GIGA_SYNC=1
Environment=GIGA_MIND=global
Environment=GIGA_CONSOLIDATE=1
ExecStart=/bin/bash -c 'DOPPLER_TOKEN=$(cat $CREDENTIALS_DIRECTORY/doppler-token) /usr/local/bin/doppler run --silent -- __HOME__/.venv-pocket/bin/python3 __PLUGIN_ROOT__/scripts/ops-cron-pocket-watcher.py'
StandardOutput=append:__HOME__/.claude/state/pocket/watcher-stdout.log
StandardError=append:__HOME__/.claude/state/pocket/watcher-stderr.log
