# claude-ops/scripts/systemd/pocket-executor.service
[Unit]
Description=Pocket Executor v3 — stateless Python supervisor (claude -p per task)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__HOME__
Environment=HOME=__HOME__
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=POCKET_STATE_DIR=__HOME__/.claude/state/pocket
Environment=POCKET_CLAUDE_BIN=/usr/local/bin/claude
Environment=POCKET_EXEC_CWD=__HOME__
Environment=POCKET_MAX_CONCURRENT=3
Environment=POCKET_WORKER_MODEL=claude-sonnet-4-6
Environment=POCKET_WORKER_TIMEOUT=1800
ExecStart=__HOME__/.venv-pocket/bin/python3 __HOME__/claude-ops/claude-ops/scripts/ops-cron-pocket-executor.py
StandardOutput=append:__HOME__/.claude/state/pocket/executor-stdout.log
StandardError=append:__HOME__/.claude/state/pocket/executor-stderr.log
# v3 requires this: workers are claude -p subprocesses that must outlive
# the oneshot service unit (they take 30-120s to complete). Without
# KillMode=process, systemd SIGKILLs the worker process group when the
# unit transitions to inactive, leaving 0-byte stdout and a bogus
# completion receipt. process mode kills only the main PID (the python
# supervisor), letting worker children run to completion under PID 1.
KillMode=process
