# claude-ops/templates/com.claude-ops.monthly-credit-reclaim.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Monthly Anthropic credit reclaim cron (HEA-4049).

  Schedule: day 1 of each month, 09:00 local time (Europe/Amsterdam for the owner).
  Equivalent cron: 0 9 1 * *

  The installer (scripts/install-monthly-credit-reclaim.sh) substitutes
  __SCRIPT_PATH__, __LOG_DIR__, __HOME__, __NODE_BIN__, __SLACK_WEBHOOK_URL__
  (from SLACK_WEBHOOK_URL at install time; empty if unset) and writes the
  rendered copy to ~/Library/LaunchAgents/com.claude-ops.monthly-credit-reclaim.plist
  before running:
    launchctl bootstrap "gui/$(id -u)" <plist>

  To uninstall:
    launchctl bootout "gui/$(id -u)/com.claude-ops.monthly-credit-reclaim"
    rm ~/Library/LaunchAgents/com.claude-ops.monthly-credit-reclaim.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-ops.monthly-credit-reclaim</string>

  <key>ProgramArguments</key>
  <array>
    <string>__NODE_BIN__</string>
    <string>__SCRIPT_PATH__</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/monthly-credit-reclaim-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/monthly-credit-reclaim-stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>__HOME__</string>
    <key>NODE_NO_WARNINGS</key>
    <string>1</string>
    <key>SLACK_WEBHOOK_URL</key>
    <string>__SLACK_WEBHOOK_URL__</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
