# scripts/systemd/novnc.service
[Unit]
Description=noVNC websockify proxy (HTTP/WebSocket -> local VNC)
Documentation=https://github.com/Lifecycle-Innovations-Limited/claude-ops/blob/main/scripts/systemd/README-novnc-vnc.md
After=graphical-session.target network-online.target

[Service]
Type=simple
# Serves the noVNC web client and bridges WebSocket <-> the local TigerVNC
# server (TigerVNC display :N -> 5900+N). Bound to 0.0.0.0 so it is reachable over
# the Tailscale interface; the box has no public IP, so this is tailnet-only.
# Ports are substituted by setup-novnc-desktop.sh (__NOVNC_PORT__, __VNC_PORT__).
ExecStart=/usr/bin/python3 -m websockify --web %h/noVNC 0.0.0.0:__NOVNC_PORT__ localhost:__VNC_PORT__
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
