# claude-ops/scripts/whatsapp/systemd/whatsapp-transcribe.service
[Unit]
Description=WhatsApp voice-note transcriber — fills empty audio rows with Whisper transcripts (oneshot)
Documentation=https://github.com/lharries/whatsapp-mcp
After=whatsapp-bridge.service
Requires=whatsapp-bridge.service
PartOf=whatsapp-bridge.service

[Service]
Type=oneshot
# OPENAI_API_KEY (and other secrets) for the metered Whisper API.
EnvironmentFile=%h/.config/systemd/env/mcp-secrets.env
# Wait up to 10s for the bridge :8080 before downloading media.
ExecStartPre=/bin/sh -c 'for i in 1 2 3 4 5; do ss -ltn 2>/dev/null | grep -q ":8080" && exit 0; sleep 2; done; exit 0'
ExecStart=/usr/bin/python3 %h/.local/share/whatsapp-mcp/transcriber/transcribe_voice_notes.py --days 2 --max 80
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whatsapp-transcribe
