# claude-ops/scripts/whatsapp/systemd/whatsapp-bridge-keepalive.service
[Unit]
Description=WhatsApp bridge keepalive (hang-detection so phone-sends always sync)
After=whatsapp-bridge.service

[Service]
Type=oneshot
ExecStart=%h/bin/wa-bridge-keepalive.sh
