# claude-ops/email-cos/systemd/email-cos-tg-process.service
[Unit]
Description=email-cos Telegram approval-callback processor
[Service]
Type=oneshot
ExecStart=%h/.local/share/email-cos/email-cos-tg-process.sh
TimeoutStartSec=120
