# claude-ops/email-cos/systemd/email-cos-sender.service
[Unit]
Description=email-cos approved-reply sender (gog, hook-exempt)

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/email-cos/config.sh
ExecStart=%h/.local/share/email-cos/email-cos-sender.sh
TimeoutStartSec=180
