# claude-ops/email-cos/systemd/email-cos-compact.service
[Unit]
Description=email-cos log/ledger compaction

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/email-cos/config.sh
ExecStart=%h/.local/share/email-cos/email-cos-compact.sh
