# claude-ops/email-cos/systemd/email-cos-approve.service
[Unit]
Description=email-cos approval reply reader (Gmail + Slack + Telegram)

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/email-cos/config.sh
ExecStart=/usr/bin/python3 %h/.local/share/email-cos/email-cos-approve-agent.py
TimeoutStartSec=180
