# claude-ops/email-cos/systemd/pocket-responder-process.service
[Unit]
Description=Pocket/social approval responder — drain responder-inbox (taps + freeform text)
[Service]
Type=oneshot
ExecStart=%h/.local/share/email-cos/pocket-responder-run.sh process
TimeoutStartSec=120
