# claude-ops/email-cos/systemd/pocket-responder-process.timer
[Unit]
Description=Drain Pocket/social approval inbox ~every 60s
[Timer]
OnBootSec=2min
OnUnitInactiveSec=60s
AccuracySec=15s
[Install]
WantedBy=timers.target
